# files/system/usr/lib/systemd/system/taidan-initial-setup-reconfiguration.service.d/override.conf
[Unit]
ConditionPathExists=
ConditionPathExists=/var/lib/taidan

[Service]
ExecStartPost=
ExecStartPost=systemctl disable taidan-initial-setup-reconfiguration
ExecStartPost=rm /var/lib/taidan

