# files/system/etc/profile.d/blue95-fastfetch.sh
alias fastfetch='fastfetch -c /usr/share/ublue-os/fastfetch.jsonc'
alias neofetch='fastfetch -c /usr/share/ublue-os/fastfetch.jsonc'
